# packages/root/etc/udev/rules.d/99-writeintelbacklight.rules
ACTION=="change", SUBSYSTEM=="backlight", RUN+="/usr/sbin/writeintelbacklight.sh"